# b2b-api.service
[Unit]
Description=B2B API for EGE Superbot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ege-bot
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 -m uvicorn b2b_api.app:app --host 127.0.0.1 --port 8001 --workers 4
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=b2b-api

[Install]
WantedBy=multi-user.target
